# deploy/configs/supervisor-nxt8.conf
# =====================================================================
# /etc/supervisor/conf.d/nxt8.conf
# Запуск FastAPI backend через uvicorn под управлением supervisor.
# =====================================================================

[program:nxt8-backend]
command=/home/nxt8/app/backend/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8001 --workers 2 --proxy-headers
directory=/home/nxt8/app/backend
user=nxt8
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=20
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/nxt8/backend.err.log
stdout_logfile=/var/log/nxt8/backend.out.log
stderr_logfile_maxbytes=20MB
stdout_logfile_maxbytes=20MB
stderr_logfile_backups=5
stdout_logfile_backups=5
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"

# Frontend строится в статику (yarn build) и отдаётся nginx — поэтому отдельная
# программа frontend здесь не нужна. Это экономит ~250 MB RAM на сервере.
